# platform_scripts/nasos/source/rc.local
#!/bin/bash

/etc/init.d/remote-mediaserver $1

exit 0
